# luci/luci-app-xl2tpd/root/etc/uci-defaults/xl2tpd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xl2tpd[-1]
	add ucitrack xl2tpd
	set ucitrack.@xl2tpd[-1].init=xl2tpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
