# packages/pollmydevice/files/pollmydevice.init
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

DIR=/sbin
PROG=pollmydevice

error()
{
	echo "${initscript}:" "$@" 1>&2
}

check_client_id()
{
	[ ! -d "/etc/pollmydevice" ] && mkdir -p "/etc/pollmydevice"
	[ ! -d "/tmp/pollmydevice" ] && mkdir -p "/tmp/pollmydevice"
	INTERFACE=$1
	if [[ -e "/etc/pollmydevice/client_id$INTERFACE" ]]; then
		cp /etc/pollmydevice/client_id$INTERFACE /tmp/pollmydevice/client_id$INTERFACE
	else
		. /usr/share/libubox/jshn.sh
		json_load "$(ubus call system board)"
		json_get_var MODEL model
		case "$MODEL" in
			"TELEOFIS LT70")
				CLIENT_ID=$(head -c 1000 /dev/urandom | tr -dc '0-9' | cut -c -15)
			;;
			"TELEOFIS RTUx68")
				SERIALNUM=$((0x$(getserialnum)))
				CLIENT_ID=$((SERIALNUM*100000+INTERFACE))
			;;
			"TELEOFIS RTUx68 RE-DE")
				SERIALNUM=$((0x$(getserialnum)))
				CLIENT_ID=$((SERIALNUM*100000+INTERFACE))
				if [ -e /etc/pollmydevice/exp_port ]; then
					/etc/pollmydevice/exp_port
					rm /etc/pollmydevice/exp_port
				fi
			;;
		esac
		echo $CLIENT_ID > /tmp/pollmydevice/client_id$INTERFACE
	fi
}

start_service_daemon()
{
	procd_open_instance
	procd_set_param command "$DIR/$PROG"
	procd_set_param respawn
	procd_close_instance
}

start_service()
{
	config_load 'pollmydevice'
	config_foreach check_client_id 'interface'
	start_service_daemon
}

service_triggers()
{
	procd_add_reload_trigger 'pollmydevice'
}

stop()
{
	service_stop $DIR/$PROG
}

# reload_service()
# {
# 	service_stop $DIR/$PROG
# 	start_service
# }

reload_service() {
        stop
        start
}