# packages/base-files-rtux68/files/etc/uci-defaults/40_feeds
#!/bin/sh

cp /etc/opkg/distfeeds_cust.conf /etc/opkg/distfeeds.conf
rm /etc/opkg/distfeeds_cust.conf

exit 0
