# luci/luci-app-hotplug/root/etc/uci-defaults/hotplug
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@hotplug[-1]
	add ucitrack hotplug
	set ucitrack.@hotplug[-1].init=hotplug
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
