# luci/luci-app-zabbix/root/etc/uci-defaults/zabbix
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@zabbix[-1]
	add ucitrack zabbix
	set ucitrack.@zabbix[-1].init=zabbix_agentd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
